# files/etc/crontabs/root
0 */6 * * * [ -x /etc/init.d/dnscrypt-proxy ] && /etc/init.d/dnscrypt-proxy restart > /dev/null 2>&1 &
0 */6 * * * [ -x /etc/init.d/miniupnpd ] && /etc/init.d/miniupnpd restart > /dev/null 2>&1 &
0 1 * * * /etc/chinadns-ng/updateChinadns-ngConfig.sh > /dev/null 2>&1 &
0 7 * * * /etc/updateGEOConfig.sh > /dev/null 2>&1 &
